# web/slick_logger/src/supervisord.conf
[supervisord]
nodaemon=true

[program:httpd]
command=/usr/local/bin/httpd-foreground

[program:tcpdump]
command=/usr/local/bin/tcpdump.sh