# web/note/src/nginx/supervisord.conf
[supervisord]
nodaemon=true

[program:nginx]
command=/usr/sbin/nginx -g "daemon off;"

[program:tcpdump]
command=/usr/local/bin/tcpdump.sh