# web/beginners_web/src/supervisord.conf
[supervisord]
nodaemon=true

[program:app]
directory=/root/app
command=/usr/local/bin/npm start

[program:tcpdump]
command=/usr/local/bin/tcpdump.sh